# contrib/systemd/engarde.service.d/bind-device.conf
[Service]
# Linux client mode binds TCP carrier sockets to interfaces with SO_BINDTODEVICE.
# Use this drop-in only if your deployment needs that behavior and the default
# service fails to create interface-bound sockets with a permission error.
# This grants CAP_NET_RAW only; the base service runs without capabilities.
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
